# crontab
# 每天 UTC 01:00（北京时间 09:00）执行证书检查
0 1 * * * root /usr/local/bin/python /app/main.py >> /var/log/ssl-monitor.log 2>&1
